# crontab
*/5 * * * * root cd /usr/src/repository && /usr/src/git-sync >> /var/log/cron.log 2>&1
