# crontab
*/5 * * * * root cd "$REPO_DIR" && /usr/src/git-sync >> /var/log/cron.log 2>&1
